# sparky.service
[Unit]
Description=Sparky

[Service]
ExecStart=/opt/sparky/sparky
WorkingDirectory=/opt/sparky
Restart=always

[Install]
WantedBy=multi-user.target
